# collections/ansible_collections/lv3/platform/roles/control_plane_recovery/templates/lv3-control-plane-backup.service.j2
[Unit]
Description=LV3 control-plane backup export
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart={{ control_plane_recovery_runtime_backup_script }}
