# collections/ansible_collections/lv3/platform/roles/cert_renewal_timer/templates/cert-renew.service.j2
[Unit]
Description={{ cert_renewal_timer_description }}
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={{ cert_renewal_timer_user }}
Group={{ cert_renewal_timer_group }}
ExecStart={{ cert_renewal_timer_script_path }}
