# collections/ansible_collections/lv3/platform/roles/server_resident_reconciliation/templates/lv3-server-resident-reconciliation.timer.j2
[Unit]
Description=Run LV3 server-resident ansible-pull reconciliation

[Timer]
OnCalendar={{ server_resident_reconciliation_on_calendar }}
RandomizedDelaySec={{ server_resident_reconciliation_randomized_delay_sec }}
Persistent=true

[Install]
WantedBy=timers.target
