# collections/ansible_collections/lv3/platform/roles/common/templates/openbao-agent-systemd-credentials.service.j2
[Unit]
Description=OpenBao Agent systemd credential delivery for {{ common_openbao_systemd_credentials_service_name }}
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=-/usr/bin/docker rm -f {{ common_openbao_systemd_credentials_container_name }}
ExecStart=/usr/bin/docker run --rm --name {{ common_openbao_systemd_credentials_container_name }} --network host --user 0:0 --cap-drop ALL --security-opt no-new-privileges:true --read-only --tmpfs /tmp:mode=1777,exec --entrypoint {{ common_openbao_systemd_credentials_container_entrypoint }} -v {{ common_openbao_systemd_credentials_agent_dir }}:/openbao-agent:ro -v {{ common_openbao_systemd_credentials_credential_dir }}:{{ common_openbao_systemd_credentials_credential_dir }}:rw {{ common_openbao_systemd_credentials_image }} agent -config=/openbao-agent/agent.hcl
ExecStartPost=/bin/bash -lc 'for attempt in $(seq 1 30); do{% if common_openbao_systemd_credentials_token_file | default("") | length > 0 %} if [ -s "{{ common_openbao_systemd_credentials_token_file }}" ] &&{% endif %} [ -s "{{ common_openbao_systemd_credentials_credential_file }}" ]; then exit 0; fi; sleep 1; done; echo "OpenBao Agent did not render systemd credentials for {{ common_openbao_systemd_credentials_service_name }} in time" >&2; exit 1'
ExecStop=-/usr/bin/docker stop {{ common_openbao_systemd_credentials_container_name }}
Restart=always
RestartSec=5s
TimeoutStartSec=60
TimeoutStopSec=20
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
