# collections/ansible_collections/lv3/platform/roles/repo_deploy_image_cache/templates/lv3-repo-deploy-image-cache.timer.j2
[Unit]
Description=Refresh approved repo-deploy base images on the Coolify deployment lane

[Timer]
OnCalendar={{ repo_deploy_image_cache_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
