# collections/ansible_collections/lv3/platform/roles/identity_core_watchdog/templates/lv3-identity-watchdog.service.j2
[Unit]
Description={{ identity_watchdog_description }}
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart={{ identity_watchdog_script_path }}
TimeoutStartSec=30
Restart=no
SyslogIdentifier={{ identity_watchdog_syslog_identifier }}
