# collections/ansible_collections/lv3/platform/roles/restic_config_backup/templates/lv3-restic-config-backup.service.j2
[Unit]
Description=LV3 restic config backup
After=network-online.target {{ restic_config_backup_runtime_openbao_agent_service_name }}.service
Wants=network-online.target
Requires={{ restic_config_backup_runtime_openbao_agent_service_name }}.service

[Service]
Type=oneshot
User=root
Group=root
UMask=0077
WorkingDirectory={{ restic_config_backup_repo_checkout_host_path }}
ExecStart={{ restic_config_backup_runtime_script }}
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{ restic_config_backup_runtime_state_dir }} {{ restic_config_backup_repo_checkout_host_path }}/receipts
