# collections/ansible_collections/lv3/platform/roles/netbox_runtime/templates/netbox-housekeeping.service.j2
[Unit]
Description=Run NetBox housekeeping
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart={{ netbox_housekeeping_script }}
