# collections/ansible_collections/lv3/platform/roles/netbox_runtime/templates/netbox-housekeeping.timer.j2
[Unit]
Description=Schedule NetBox housekeeping

[Timer]
OnCalendar={{ netbox_housekeeping_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
