# collections/ansible_collections/lv3/platform/roles/data_retention/templates/lv3-data-retention.timer.j2
[Unit]
Description=Run LV3 receipt and mutation-audit retention

[Timer]
OnCalendar={{ data_retention_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
